# bash/pia-router/etc/systemd/system/monitor-eth1.service
[Unit]
Description=Monitor changes to eth1 and trigger a script
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
Restart=always
StartLimitBurst=10
StartLimitIntervalSec=1
ExecStart=/opt/script/monitor-eth1
User=root
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
