# bash/pia-router/etc/systemd/system/piaqbittorrentport.service
[Unit]
Description=Update qbittorrent port from PIA

[Service]
Type=oneshot
ExecStart=usr/bin/env bash /opt/script/piaqbittorrentport